# deploy/server/queue-app.service
# /etc/systemd/system/queue-app.service
# Podman-managed Queue App container. Pulls the latest image before each start
# while falling back to the cached copy if the host is offline.

[Unit]
Description=Queue App Container (Podman)
Wants=network-online.target
After=network-online.target

[Service]
EnvironmentFile=-/etc/default/queue-app

Restart=always
RestartSec=10

# Ensure the container uses the freshest image and starts cleanly.
ExecStartPre=/bin/sh -c 'touch /etc/default/queue-app'
ExecStartPre=/bin/sh -c '/usr/bin/podman pull ghcr.io/ledokun/028-simple-queue-system:latest || /usr/bin/podman image exists ghcr.io/ledokun/028-simple-queue-system:latest'
ExecStartPre=-/usr/bin/podman stop queue-app
ExecStartPre=-/usr/bin/podman rm queue-app

ExecStart=/usr/bin/podman run \
  --name queue-app \
  --publish 3000:3000 \
  --env-file /etc/default/queue-app \
  ghcr.io/ledokun/028-simple-queue-system:latest

ExecStop=-/usr/bin/podman stop --time 10 queue-app
ExecStopPost=-/usr/bin/podman rm queue-app

# Persist container logs to journald.
Type=simple
KillMode=control-group

[Install]
WantedBy=multi-user.target
