# deploy/quadlets/queue-app.container
# /etc/containers/systemd/queue-app.container

[Unit]
Description=Queue App Container
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=queue-app
Image=ghcr.io/ledokun/028-simple-queue-system:latest
PublishPort=3000:3000

# Environment Variables
Environment=RUST_LOG=info
Environment=SERVER_ADDRESS=0.0.0.0
Environment=SERVER_PORT=3000
Environment=MAX_HISTORY_SIZE=5
Environment=MAX_SKIPPED_HISTORY_SIZE=5
Environment=SERVE_DIR_PATH=./public
Environment=ANNOUNCEMENTS_AUDIO_SUB_PATH=media/announcements
Environment=BANNERS_SUB_PATH=media/banners
Environment=ANNOUNCEMENT_AUTO_CYCLE_INTERVAL_SECONDS=1200
Environment=ANNOUNCEMENT_MANUAL_TRIGGER_COOLDOWN_SECONDS=5
Environment=BANNER_ROTATION_INTERVAL_SECONDS=10
Environment=GTTS_CACHE_BASE_PATH=/tmp/gtts_audio_cache
Environment=TTS_CACHE_MAXIMUM_FILES=500
Environment=TTS_EXTERNAL_SERVICE_TIMEOUT_SECONDS=15
Environment=TTS_SUPPORTED_LANGUAGES="th:Thai,en-uk:British English"
Environment=SSE_KEEP_ALIVE_INTERVAL_SECONDS=15
Environment=SSE_EVENT_BUFFER_SIZE=200
Environment=TTS_CACHE_WEB_PATH=/tts_cache

[Service]
Restart=always

[Install]
WantedBy=multi-user.target
